# substandard/ruptime.timer
[Unit]
Description=ruptime

[Timer]
OnBootSec=1m

[Install]
WantedBy=basic.target
