# systemd/ruptime.service
[Unit]
Description=Update ruptime entry

[Service]
Type=oneshot
ExecStart=/usr/bin/ruptime -u
